# com.shopmanager.autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shopmanager.autostart</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ruta/completa/al/proyecto/start_silent.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/shopmanager-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/shopmanager-stderr.log</string>
</dict>
</plist>
